# deployment/services.d/werkhof-ichen-backend.ini
[program:werkhof-ichen-backend]
directory=/home/josen/werkhof-ichen/backend
command=/usr/bin/npm run start
autostart=yes
autorestart=yes
