# deployment/services.d/werkhof-ichen-webhook.ini
[program:werkhof-ichen-webhook]
command=/home/josen/bin/webhook -hooks /home/josen/werkhof-ichen/deployment/hooks.json -verbose
autostart=yes
autorestart=yes
